# Lesson_7_Systemd/files/logmon.service
[Unit]
Description=Log test with timer

[Service]
Type=simple
ExecStart=/bin/bash /vagrant/log_count.sh

[Install]
WantedBy=multi-user.target
