# Lesson_7_Systemd/logmon.service
[Unit]
Description=Log test with timer

[Service]
Type=simple
ExecStart=/bin/bash /vagrant/script.sh

[Install]
WantedBy=multi-user.target
